# fanctl.service
[Unit]
Description=Fan Control

[Service]
Type=simple
ExecStart=/usr/sbin/fanctl >/var/log/fanctl.log

